# foswvs.service
[Unit]
Description=foswvs

[Service]
ExecStart=/home/pi/foswvs/setup.sh

[Install]
WantedBy=multi-user.target
